# systemd/discord-component-v2-broker.service.template
[Unit]
Description=discord-component-v2 broker (Discord Gateway)
After=network-online.target

[Service]
Type=simple
WorkingDirectory=__INSTALL_DIR__
Environment=PYTHONUNBUFFERED=1
Environment=DISCORD_COMPONENT_V2_WORKSPACE_ID=__WORKSPACE_ID__
ExecStart=__PYTHON__ __INSTALL_DIR__/scripts/broker_gateway.py
Restart=always
RestartSec=2

[Install]
WantedBy=default.target
