# systemd/discord-component-v2-worker.service.template
[Unit]
Description=discord-component-v2 worker (normalize + inject)
After=network-online.target

[Service]
Type=simple
WorkingDirectory=__INSTALL_DIR__
Environment=PYTHONUNBUFFERED=1
ExecStart=__PYTHON__ __INSTALL_DIR__/scripts/worker.py
Restart=always
RestartSec=2

[Install]
WantedBy=default.target
